# home/dot_config/systemd/user/backup-packages.service
[Unit]
Description=Backup Packages

[Service]
Type=oneshot
ExecStart=%h/.local/bin/chezmoi-externals-backup