# home/dot_config/systemd/user/backup-packages.timer
[Unit]
Description=Run backup-packages daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
