# home/dot_config/systemd/user/codex-lb.service
# ~/.config/systemd/user/codex-lb.service
[Unit]
Description=Codex Load Balancer (uvx)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/bin/sh -c 'exec uvx codex-lb'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
