# home/dot_config/systemd/user/dim-screen.service
[Unit]
Description=Dim the screen via hyprsunset

[Service]
Type=oneshot
ExecStart=/usr/bin/hyprsunset --gamma 0.1
