# home/dot_config/systemd/user/dim-screen.timer
[Unit]
Description=Timer to dim the screen after a period of inactivity

[Timer]
OnActiveSec=15s

[Install]
WantedBy=timers.target
